# systemd/user/cycle-wallpaper.service
[Unit]
Description=Cycle sway wallpaper

[Service]
Type=oneshot
ExecStart=%h/git/nordfjord/dotfiles/bin/cycle-wallpaper
